# crontab
30 15 * * * certbot renew --deploy-hook "nginx -s reload"
# Mandatory empty line.
