# crontab
* * * * * certbot renew --dry-run --deploy-hook "nginx -s reload" && echo "Cron job started." >> ~/logs.txt
# Mandatory empty line.
